# forum-flask.conf
[program:forum-flask]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000
directory=/var/www/forum-flask
autostart=true
autorestart=true
